# home-dashboard.socket
[Unit]
Description=Home Dashboard Socket (port 80)

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=sockets.target
